# elek-board-viewer/systemd/boardview-cleanup.timer
[Unit]
Description=Daily cleanup for boardview screenshots and logs

[Timer]
OnCalendar=daily
Persistent=true
Unit=boardview-cleanup.service

[Install]
WantedBy=timers.target
