# elek-board-viewer/systemd/flexbv-headless.service
[Unit]
Description=FlexBV headless Xvfb/Openbox session
After=default.target

[Service]
Type=simple
ExecStart=/home/ricky/builds/elek-board-viewer/scripts/flexbv_session.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
