# agents/_shared/bridge/systemd/claude-backup@.service
[Unit]
Description=Claude agents nightly restic snapshot — %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ricky
Group=ricky
Environment=TZ=Asia/Singapore
EnvironmentFile=/home/ricky/builds/agents/_shared/bridge/envs/%i.env
ExecStart=/usr/bin/restic backup /home/ricky/builds/agents/_shared/state/ --verbose
StandardOutput=append:/home/ricky/builds/agents/_shared/logs/restic-backup.log
StandardError=append:/home/ricky/builds/agents/_shared/logs/restic-backup.log
