# agents/_shared/bridge/systemd/claude-backup@.timer
[Unit]
Description=Nightly restic timer — %i

[Timer]
OnCalendar=*-*-* 03:00:00 Asia/Singapore
Persistent=true
Unit=claude-backup@%i.service

[Install]
WantedBy=timers.target
