# backmarket/services/bm-payout/bm-payout.service
[Unit]
Description=BM Payout Webhook Service (SOP 03b)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ricky/builds/backmarket/services/bm-payout
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/ricky/config/.env

[Install]
WantedBy=default.target
