# monday/services/repair-stock-check/repair-stock-check.service
[Unit]
Description=Monday Repair Stock Check Service (Requested Repairs -> Monday update)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ricky/builds/monday/services/repair-stock-check
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/ricky/config/.env
Environment=SHADOW_MODE=false
Environment=PORT=8015

[Install]
WantedBy=default.target
